# testdbus-1.0/data/testdbus.service.in
[Unit]
Description=sercure box service

[Service]
Type=dbus
BusName=org.freedesktop.Testdbus
ExecStart=@testdbusprivdir@/testdbusd
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target